# archive/systemd/service-health-check.timer
[Unit]
Description=Regular service health check timer
Documentation=https://docs.nova.local/monitoring

[Timer]
# Run every minute
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=1s

[Install]
WantedBy=timers.target